# omnibus/config/templates/datadog-agent/systemd.logs.service.erb
[Unit]
Description="Datadog Agent"
Requisite=datadog-agent.service
After=network.target datadog-agent.target
PartOf=datadog-agent.service

[Service]
Type=simple
PIDFile=<%= install_dir %>/run/logs-agent.pid
User=dd-agent
Restart=on-failure
ExecStart=<%= install_dir %>/bin/agent/logs-agent --ddconfd=<%= etc_dir %>/conf.d --ddconfig=<%= etc_dir %> --pid=<%= install_dir %>/run/logs-agent.pid

[Install]
WantedBy=multi-user.target
