# omnibus/config/templates/datadog-agent/upstart_debian.trace.conf.erb
description "Datadog Agent"

start on (started networking and started datadog-agent)
stop on stopping datadog-agent

respawn

# Logging to console from the agent is disabled since the agent already logs using file or
# syslog depending on its configuration. We make upstart log what the process still outputs in order
# to log panics/crashes to /var/log/upstart/datadog-agent.log
console log
env DD_LOG_TO_CONSOLE=false

setuid dd-agent

script
  exec <%= install_dir %>/bin/agent/trace-agent --config <%= etc_dir %>/trace_agent.conf --pid <%= install_dir %>/run/trace-agent.pid
end script

post-stop script
  rm -f <%= install_dir %>/run/trace-agent.pid
end script
